# package/preremove.sh
#!/bin/sh

systemctl disable mqtt-homekit-light@
